# scripts/api/nix-for-humanity-api.service
[Unit]
Description=Nix for Humanity REST API Server
Documentation=https://github.com/Luminous-Dynamics/nix-for-humanity
After=network.target

[Service]
Type=simple
User=nix-api
Group=nix-api
WorkingDirectory=/srv/luminous-dynamics/11-meta-consciousness/nix-for-humanity
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="PYTHONPATH=/srv/luminous-dynamics/11-meta-consciousness/nix-for-humanity/scripts"
Environment="API_HOST=0.0.0.0"
Environment="API_PORT=5000"
Environment="API_DEBUG=false"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/luminous-dynamics/11-meta-consciousness/nix-for-humanity/data

# Start command
ExecStart=/usr/bin/python3 scripts/api/nix_api_server.py

# Restart policy
Restart=always
RestartSec=5

# Resource limits
LimitNOFILE=4096
MemoryLimit=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nix-api

[Install]
WantedBy=multi-user.target