# scripts/nix-for-humanity-headless.service
[Unit]
Description=Nix for Humanity Headless Engine
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /srv/luminous-dynamics/11-meta-consciousness/nix-for-humanity/scripts/run_headless_server.py
Restart=on-failure
RestartSec=5

# Run as regular user, not root
User=nobody
Group=nogroup

# Security hardening
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes

# Socket permissions
RuntimeDirectory=nix-for-humanity
RuntimeDirectoryMode=0755

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="NIX_FOR_HUMANITY_DATA=/var/lib/nix-for-humanity"

[Install]
WantedBy=multi-user.target
